# autocat3.service
[Unit]
Description=autocat3 Service
After=network.target

[Service]
User=autocat
Type=simple
WorkingDirectory=/var/lib/autocat/autocat3
ExecStart=/opt/rh/rh-python36/root/usr/bin/pipenv run python CherryPyApp.py

[Install]
WantedBy=multi-user.target