# autocat3.service
[Unit]
Description=autocat3 Service
After=network.target

[Service]
User=autocat
Type=simple
RuntimeDirectory=autocat
WorkingDirectory=/var/lib/autocat/autocat3
ExecStartPre=-/usr/bin/mkdir /var/run/autocat
ExecStart=/var/lib/autocat/.local/bin/pipenv run python CherryPyApp.py

[Install]
WantedBy=multi-user.target
